# debian/kioskme.xorg.service
[Unit]
Description=xorg

[Service]
Type=simple
Restart=on-failure
ExecStart=/usr/bin/X
